# conf/supervisord/supervisord.conf
[supervisord]

[supervisorctl]

[program:ft-uwsgi]
command = /usr/local/bin/uwsgi --ini /ft/conf/uwsgi/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
